# etc/crontab
# m h  dom mon dow   command
0 5  * * *  /root/home/www/code/test/backend-django/cron.sh