# etc/crontab
# m h  dom mon dow   command
0 5  * * *  /root/backend-django/cron_job.sh