# etc/systemd/dev-test.socket
[Unit]
Description=dev-test socket
[Socket]
ListenStream=/run/dev-test.sock
[Install]
WantedBy=sockets.target